# caddy.service
[Unit]
Description=Caddy HTTP/2 web server
After=network.target

[Service]
User=pi
ExecStart=/usr/local/bin/caddy run --config /etc/caddy/Caddyfile
ExecReload=/usr/local/bin/caddy reload --config /etc/caddy/Caddyfile

[Install]
WantedBy=multi-user.target